# config/crontab
# Cron job definitions for Kamal cron container.
#
# This file is copied into the Kamal image and loaded by:
#   cat config/crontab | crontab -
#
# Notes:
# - Cron strips most container env vars. `BASH_ENV` reloads the required ones
#   from a startup snapshot before each job starts.
# - `CRON_TZ` keeps the schedule pinned to America/Chicago, including DST.
# - The bridge sync delta crawl uses the previous America/Chicago day at
#   00:00:00 and converts that to UTC for the bridge API.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
BASH_ENV=/app/scripts/cron_env.sh
CRON_TZ=America/Chicago
TZ=America/Chicago

# Daily 2:00 AM: enqueue bridge sync delta crawl (yesterday 00:00:00Z).
# Use venv Python: cron runs with minimal PATH and system python3 lacks requests/etc.
0 2 * * * /opt/venv/bin/python3 /app/scripts/trigger_bridge_sync_cron.py >> /proc/1/fd/1 2>> /proc/1/fd/2

# Daily 3:00 AM: enqueue blog sync (Celery)
0 3 * * * /opt/venv/bin/python3 /app/scripts/trigger_blog_sync_cron.py >> /proc/1/fd/1 2>> /proc/1/fd/2

# Daily 4:15 AM: regenerate sitemap XML after overnight content refresh jobs.
15 4 * * * /opt/venv/bin/python3 /app/scripts/generate_sitemap.py >> /proc/1/fd/1 2>> /proc/1/fd/2

# Daily 4:45 AM: roll up analytics, create future partitions, and drop expired raw partitions.
45 4 * * * /opt/venv/bin/python3 /app/backend/scripts/manage_analytics_storage.py --mode maintenance >> /proc/1/fd/1 2>> /proc/1/fd/2
